# system_files/usr/lib/systemd/system/usr-local.mount
[Unit]
DefaultDependencies=no
Before=local-fs.target sysinit.target

[Mount]
What=/var/usrlocal
Where=/usr/local
Options=bind

[Install]
WantedBy=local-fs.target